# build/docker/supervisord.conf
[supervisord]
nodaemon=true
[supervisorctl]
[program:xcp]
directory=/app
command=/app/xcpserver
user=root
autostart=true
autorestart=true   
stderr_logfile=/app/xcp.err.log
stdout_logfile=/app/xcp.out.log


